# fleet-manager.service
[Unit]
Description=Bitcoin Mining Fleet Manager (DirtySats)
After=network.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__INSTALL_DIR__
Environment="PATH=__INSTALL_DIR__/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=__INSTALL_DIR__/venv/bin/gunicorn -w 1 --threads 4 -b 0.0.0.0:5001 --timeout 120 app:app
Restart=always
RestartSec=10
StandardOutput=append:__INSTALL_DIR__/logs/fleet-manager.log
StandardError=append:__INSTALL_DIR__/logs/fleet-manager.error.log

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
